# examples/systemd/hermes-napcat-qq-bridge.service
[Unit]
Description=Hermes NapCat QQ Bridge
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/home/YOUR_USER
ExecStart=/home/YOUR_USER/.local/bin/hermes napcat-qq-bridge run \
  --onebot-url http://127.0.0.1:3000 \
  --onebot-token YOUR_ONEBOT_TOKEN \
  --listen-host 127.0.0.1 \
  --listen-port 8096 \
  --webhook-path /napcat \
  --allow-user YOUR_QQ \
  --allow-group YOUR_GROUP \
  --hermes-workdir /home/YOUR_USER \
  --hermes-toolsets terminal,file,web
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target

