# examples/systemd/hermes-napcat-qq-bridge.service
[Unit]
Description=Hermes NapCat QQ Personal Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h
Environment=PYTHONUNBUFFERED=1
Environment=HERMES_HOME=%h/.hermes
ExecStart=%h/.hermes/hermes-agent/venv/bin/python3 %h/.hermes/plugins/napcat_qq_bridge/bridge.py run --config-file %h/.hermes/napcat_qq_bridge/config.json
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
